# deployment/systemd/incontrol-worker.service
[Unit]
Description=InControl Celery Worker
Documentation=https://github.com/yourusername/incontrol
After=network.target mariadb.service redis.service
Requires=redis.service mariadb.service

[Service]
User=incontrol
Group=incontrol
WorkingDirectory=/opt/incontrol
ExecStart=/opt/incontrol/venv/bin/celery -A incontrol worker \
    --loglevel=info \
    --concurrency=4 \
    --max-tasks-per-child=1000 \
    --time-limit=3600 \
    --soft-time-limit=3300

# Security
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/opt/incontrol
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

Restart=always
RestartSec=5
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target