# deployment/systemd/incontrol-ws.service
[Unit]
Description=InControl WebSocket Service
Documentation=https://github.com/yourusername/incontrol
After=network.target mariadb.service redis.service
Requires=redis.service mariadb.service

[Service]
User=incontrol
Group=incontrol
WorkingDirectory=/opt/incontrol
ExecStart=/opt/incontrol/venv/bin/daphne -b 127.0.0.1 -p 8001 incontrol.asgi:application
Restart=always
RestartSec=5
TimeoutStopSec=30s

# Security
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/opt/incontrol
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

[Install]
WantedBy=multi-user.target 